# files/system-cinnamon/etc/systemd/system/lightdm.service.d/10-wait-flatpak-setup.conf
[Unit]
Wants=system-flatpak-setup.service
After=system-flatpak-setup.service
