# files/system-cinnamon/etc/systemd/system/system-flatpak-setup.service.d/10-firstboot-only.conf
[Unit]
ConditionPathExists=!/var/lib/cinnamon-flatpak-setup/done

[Service]
StateDirectory=cinnamon-flatpak-setup
ExecStartPost=/usr/bin/touch /var/lib/cinnamon-flatpak-setup/done
